# etc/network/if-up.d/ics
#!/bin/sh
#--------------------ics
. /etc/init.d/ics_sharing.sh
#--------------------ics /

