# etc/network/if-up.d/ics
#!/bin/sh
iw dev wlan0 set 4addr on
#--------------------ics
. /etc/init.d/ics_sharing.sh
#--------------------ics /

